# packaging/ubuntu/deb/boundless-server-tomcat9/var/debian/postinst
#!/bin/bash
# alien added permissions fixup code
chown 'tomcat9:tomcat9' '/etc/init.d/tomcat9'
chown 'tomcat9:tomcat9' '/etc/logrotate.d/tomcat9'
chown 'tomcat9:tomcat9' '/etc/tomcat9'
chown 'tomcat9:tomcat9' '/etc/tomcat9/Catalina'
chown 'tomcat9:tomcat9' '/etc/tomcat9/Catalina/localhost'
chown 'tomcat9:tomcat9' '/etc/tomcat9/catalina.policy'
chown 'tomcat9:tomcat9' '/etc/tomcat9/catalina.properties'
chown 'tomcat9:tomcat9' '/etc/tomcat9/context.xml'
chown 'tomcat9:tomcat9' '/etc/tomcat9/jaspic-providers.xml'
chown 'tomcat9:tomcat9' '/etc/tomcat9/jaspic-providers.xsd'
chown 'tomcat9:tomcat9' '/etc/tomcat9/java-libs'
chown 'tomcat9:tomcat9' '/etc/tomcat9/logging.properties'
chown 'tomcat9:tomcat9' '/etc/tomcat9/server-opts'
chown 'tomcat9:tomcat9' '/etc/tomcat9/server-opts/geoserver'
chown 'tomcat9:tomcat9' '/etc/tomcat9/server.xml'
chown 'tomcat9:tomcat9' '/etc/tomcat9/tomcat-users.xml'
chown 'tomcat9:tomcat9' '/etc/tomcat9/tomcat-users.xsd'
chown 'tomcat9:tomcat9' '/etc/tomcat9/tomcat9.conf'
chown 'tomcat9:tomcat9' '/etc/tomcat9/web.xml'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/bootstrap.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/catalina-tasks.xml'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/catalina.bat'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/catalina.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/ciphers.bat'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/ciphers.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/commons-daemon-native.tar.gz'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/commons-daemon.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/configtest.bat'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/configtest.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/daemon.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/digest.bat'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/digest.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/makebase.bat'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/makebase.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/setclasspath.bat'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/setclasspath.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/shutdown.bat'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/shutdown.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/startup.bat'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/startup.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/tomcat-juli.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/tomcat-native.tar.gz'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/tool-wrapper.bat'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/tool-wrapper.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/version.bat'
chown 'tomcat9:tomcat9' '/opt/tomcat9/bin/version.sh'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/annotations-api.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/catalina-ant.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/catalina-ha.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/catalina-storeconfig.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/catalina-tribes.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/catalina.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/ecj-4.9.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/el-api.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/jasper-el.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/jasper.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/jaspic-api.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/jsp-api.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/servlet-api.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-api.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-coyote.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-dbcp.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-i18n-de.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-i18n-es.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-i18n-fr.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-i18n-ja.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-i18n-ko.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-i18n-pt-BR.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-i18n-ru.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-i18n-zh-CN.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-jdbc.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-jni.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-util-scan.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-util.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/tomcat-websocket.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/lib/websocket-api.jar'
chown 'tomcat9:tomcat9' '/opt/tomcat9/webapps'
chown 'tomcat9:tomcat9' '/usr/sbin/tomcat9'
chown 'tomcat9:tomcat9' '/var/cache/tomcat9'
chown 'tomcat9:tomcat9' '/var/cache/tomcat9/temp'
chown 'tomcat9:tomcat9' '/var/cache/tomcat9/work'
chown 'tomcat9:tomcat9' '/var/log/tomcat9'

chmod +x /etc/init.d/tomcat9
chmod +x /usr/sbin/tomcat9

if [ -f /sbin/chkconfig ]; then
  chkconfig --add tomcat9
elif [ -f /usr/sbin/update-rc.d ]; then
  if [ ! -f /etc/init.d/tomcat9 ]; then
    ln -s /etc/rc.d/init.d/tomcat9 /etc/init.d/tomcat9
  fi
update-rc.d tomcat9 defaults
fi
